# systemd/user/hexe-supervisor-api.service.in
[Unit]
Description=Hexe Supervisor API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=@INSTALL_DIR@/backend
Environment=HEXE_SUPERVISOR_TRANSPORT=socket
Environment=HEXE_SUPERVISOR_SOCKET=/run/hexe/supervisor.sock
Environment=PYTHONUNBUFFERED=1
ExecStart=@INSTALL_DIR@/backend/.venv/bin/python -m app.supervisor.server
Restart=always
RestartSec=2
EnvironmentFile=%h/.config/hexe/admin.env

[Install]
WantedBy=default.target
