# systemd/user/synthia-dashboard.service.in
[Unit]
Description=Hexe Dashboard Browser
After=graphical-session.target synthia-backend.service synthia-frontend-dev.service network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/dan/.Xauthority
Environment=HOME=/home/dan
Environment=CHROME_BIN=/usr/bin/google-chrome-stable
Environment=SYNTHIA_DASHBOARD_URL=http://127.0.0.1/

ExecStartPre=/usr/bin/sleep 5
ExecStartPre=/usr/bin/bash -lc 'DISPLAY=:0 xset s off -dpms s noblank || true'
ExecStart=/usr/bin/bash -lc '"${CHROME_BIN}" \
  --new-window \
  --kiosk "${SYNTHIA_DASHBOARD_URL}" \
  --start-maximized \
  --disable-session-crashed-bubble \
  --disable-infobars \
  --no-first-run \
  --no-default-browser-check \
  --disable-gpu \
  --disable-software-rasterizer'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
